# rsp-fsystem/src/etc.skel/.config/autostart/ListenTopLeftCursor.desktop
[Desktop Entry]
Exec=/usr/lib/rsp/proc/listenTopLeftCursor
Hidden=true
Name=RspListenTopLeftCursor
Type=Application
Version=1.0
