# rsp-fsystem/src/etc.skel/.config/autostart/network-manager-applet.desktop
[Desktop Entry]
Exec=/usr/local/bin/nm-applet
Name=network-manager-applet
OnlyShowIn=RSp;
Type=Application
Version=1.0
X-RSp-Need-Tray=true
