# rsp-fsystem/src/etc.skel/.config/autostart/RspRspPanel.desktop
[Desktop Entry]
Exec=/usr/lib/rsp/proc/xcbRspPanel
Hidden=true
Name=RspRspPanel
Type=Application
Version=1.0
